# docs/examples/hermes-assistant-dev.service.example
[Unit]
Description=Hermes Assistant dev gateway
After=network-online.target sandboxed-sh-dev.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/sandboxed-sh/hermes-assistant-dev.env
WorkingDirectory=/var/lib/hermes-assistant-dev

# Install Hermes on the host first. The gateway process owns Telegram and loads
# ~/.hermes/config.yaml from HERMES_HOME.
ExecStart=/usr/local/bin/hermes gateway

Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
